# roles/xai_node/templates/xai-node.service.j2
#
# {{ service_name }}.service - XAI Blockchain Node
#
# Managed by Ansible - DO NOT EDIT MANUALLY
# Node: {{ node_name }} ({{ node_type }})
#

[Unit]
Description=XAI Blockchain Node ({{ node_name }})
Documentation=https://docs.xaiblockchain.com
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ daemon_user }}
Group={{ daemon_user }}
WorkingDirectory={{ xai_home | default('/home/ubuntu/xai') }}

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=XAI_NETWORK=mvp-testnet
Environment=XAI_CHAIN_ID={{ chain_id }}
Environment=XAI_DATA_DIR={{ daemon_home }}/data
Environment=XAI_GENESIS_FILE={{ xai_mvp_home | default('/home/ubuntu/xai-mvp') }}/genesis.json
Environment=XAI_LOG_LEVEL=INFO
Environment=XAI_MINING_ENABLED={{ mining | default(false) | lower }}
Environment=XAI_RATE_LIMIT_FAUCET=100
Environment=XAI_RATE_LIMIT_FAUCET_WINDOW=3600
Environment=XAI_API_AUTH_REQUIRED=0
Environment=XAI_RATE_LIMIT_ENABLED=0
Environment=XAI_FAUCET_WALLET_FILE={{ xai_mvp_home | default('/home/ubuntu/xai-mvp') }}/wallets/faucet_wallet.json
Environment=XAI_FAUCET_AMOUNT=10
Environment=XAI_ALLOW_EMPTY_MINING=true
Environment=XAI_FINALITY_ENABLED=true
Environment=XAI_FAST_MINING=1

# Main process
ExecStart={{ venv_path }}/bin/python -m xai.core.node \
  --host 0.0.0.0 \
  --port {{ rpc_port }} \
  --p2p-port {{ p2p_port }} \
  --data-dir {{ daemon_home }}/data \
{% if mining | default(false) and miner_address is defined and miner_address %}
  --miner {{ miner_address }} \
{% endif %}
{% if node_type == 'sentry' %}
  --sentry-mode \
{% endif %}
{% if peers is defined and peers %}
  --peers {{ peers | join(' ') }}
{% endif %}

# Restart policy
Restart=always
RestartSec=10
TimeoutStartSec=60
TimeoutStopSec=300

# Resource limits
LimitNOFILE=65535
LimitNPROC=65535

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ service_name }}

[Install]
WantedBy=multi-user.target
