# xai-node/templates/xai.service.j2
# Managed by Ansible - DO NOT EDIT MANUALLY
# XAI Blockchain Node systemd service

[Unit]
Description=XAI Blockchain Node
After=network-online.target
Wants=network-online.target

[Service]
User={{ daemon_user }}
Group={{ daemon_user }}
WorkingDirectory={{ daemon_home }}/src/xai
ExecStart={{ venv_path }}/bin/python3 core/node.py --miner {{ xai_miner_address | default('TXAI759460929623A1A488388CED4054734BCF1F835D') }} --port {{ xai_rpc_port | default(8545) }} --host 0.0.0.0 --data-dir {{ xai_data_dir | default(daemon_home) }} --p2p-port {{ xai_p2p_port | default(8333) }}
Restart=always
RestartSec=5
LimitNOFILE=65535

# Environment
Environment="PATH={{ venv_path }}/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
Environment="XAI_HOME={{ daemon_home }}"
Environment="XAI_NETWORK=testnet"
Environment="XAI_RPC_HOST=0.0.0.0"
Environment="XAI_RPC_PORT={{ xai_rpc_port | default(8545) }}"
Environment="XAI_P2P_PORT={{ xai_p2p_port | default(8333) }}"
Environment="XAI_WS_PORT={{ xai_ws_port | default(8765) }}"
Environment="XAI_ALLOW_EMPTY_MINING=1"
Environment="XAI_MINING_HEARTBEAT_SECONDS=30"
Environment="XAI_FAST_MINING=1"
Environment="XAI_MAX_TEST_MINING_DIFFICULTY=4"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=xai

[Install]
WantedBy=multi-user.target
