# roles/monitoring/templates/prometheus.service.j2
# Managed by Ansible - DO NOT EDIT MANUALLY
# Prometheus systemd service

[Unit]
Description=Prometheus
Wants=network-online.target
After=network-online.target

[Service]
User=prometheus
Group=prometheus
Type=simple
ExecStart=/usr/local/bin/prometheus \
    --config.file=/etc/prometheus/prometheus.yml \
    --storage.tsdb.path=/var/lib/prometheus/ \
    --web.listen-address=:{{ prometheus_port }} \
    --web.console.templates=/etc/prometheus/consoles \
    --web.console.libraries=/etc/prometheus/console_libraries \
    --storage.tsdb.retention.time=15d
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
